# server/systemd/deploy-hub-prune.timer
[Unit]
Description=deploy-hub: daily image prune

[Timer]
OnCalendar=*-*-* 04:10:00 UTC
Persistent=true

[Install]
WantedBy=timers.target
